# systemd/.config/systemd/user/ra-multiplex.service
[Unit]
Description=Rust Analyzer Multiplex

[Service]
ExecStart=/home/seb/.cargo/bin/ra-multiplex server

[Install]
WantedBy=default.target
